# recipes-devtools/wendyos-containerd-registry/files/wendyos-dev-registry.service
[Unit]
Description=WendyOS Development Container Registry
Documentation=https://github.com/wendylabsinc/containerd-registry
After=containerd.service network-online.target wendyos-dev-registry-import.service
Requires=containerd.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
Environment=REGISTRY_NAMESPACE=default
Environment=REGISTRY_NAME=wendyos-dev-registry
Environment=REGISTRY_IMAGE=wendyos/containerd-registry:v1.2.0
Environment=LISTEN_ADDRESS=0.0.0.0:5000

# Pre-flight checks: Verify image exists, import if missing (self-healing)
ExecStartPre=/bin/sh -c '\
  if ! /usr/bin/ctr -n ${REGISTRY_NAMESPACE} images ls | /bin/grep -q ${REGISTRY_IMAGE}; then \
    echo "Registry image not found, importing from archive..."; \
    /usr/bin/ctr -n ${REGISTRY_NAMESPACE} images import /usr/share/wendyos/offline-images/containerd-registry.tar; \
    /usr/bin/ctr -n ${REGISTRY_NAMESPACE} images tag wendyos/containerd-registry:v1.2.0 wendyos/containerd-registry:latest || true; \
    echo "Registry image imported successfully"; \
  fi'
# Cleanup any existing container
ExecStartPre=-/bin/sh -c '/usr/bin/ctr -n ${REGISTRY_NAMESPACE} tasks kill ${REGISTRY_NAME} 2>/dev/null; /usr/bin/ctr -n ${REGISTRY_NAMESPACE} tasks delete ${REGISTRY_NAME} 2>/dev/null; /usr/bin/ctr -n ${REGISTRY_NAMESPACE} containers delete ${REGISTRY_NAME} 2>/dev/null; true'

# Start the registry container
ExecStart=/usr/bin/ctr -n ${REGISTRY_NAMESPACE} run \
    --detach \
    --net-host \
    --mount type=bind,src=/run/containerd/containerd.sock,dst=/run/containerd/containerd.sock,options=rbind:rw \
    --env LISTEN_ADDRESS=${LISTEN_ADDRESS} \
    --env CONTAINERD_NAMESPACE=${REGISTRY_NAMESPACE} \
    --env LOG_FORMAT=json \
    ${REGISTRY_IMAGE} \
    ${REGISTRY_NAME}

# Graceful shutdown
ExecStop=/bin/sh -c '/usr/bin/ctr -n ${REGISTRY_NAMESPACE} tasks kill ${REGISTRY_NAME} 2>/dev/null || true; /usr/bin/ctr -n ${REGISTRY_NAMESPACE} tasks wait ${REGISTRY_NAME} 2>/dev/null || true'
ExecStopPost=-/bin/sh -c '/usr/bin/ctr -n ${REGISTRY_NAMESPACE} tasks delete ${REGISTRY_NAME} 2>/dev/null; /usr/bin/ctr -n ${REGISTRY_NAMESPACE} containers delete ${REGISTRY_NAME} 2>/dev/null; true'

# No automatic restart for oneshot services - handled by wendy-agent if needed
TimeoutStartSec=60s
TimeoutStopSec=45s

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wendyos-dev-registry

# Security
NoNewPrivileges=true

[Install]
# Disabled by default - started on-demand by wendy-agent
WantedBy=multi-user.target
